# deploy/macos/com.opc-studio.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.opc-studio</string>
  <key>ProgramArguments</key>
  <array>
    <string>__ROOT__/.venv/bin/uvicorn</string>
    <string>app.main:app</string>
    <string>--host</string>
    <string>__HOST__</string>
    <string>--port</string>
    <string>__PORT__</string>
  </array>
  <key>WorkingDirectory</key>
  <string>__ROOT__</string>
  <key>EnvironmentVariables</key>
  <dict>
    <key>PYTHONPATH</key>
    <string>__ROOT__/backend/app</string>
    <key>OPC_DATA_DIR</key>
    <string>__OPC_DATA_DIR__</string>
    <key>OPC_SECRET_KEY</key>
    <string>__OPC_SECRET_KEY__</string>
  </dict>
  <key>RunAtLoad</key>
  <true/>
  <key>KeepAlive</key>
  <true/>
  <key>StandardOutPath</key>
  <string>__ROOT__/data/logs/opc-studio.out.log</string>
  <key>StandardErrorPath</key>
  <string>__ROOT__/data/logs/opc-studio.err.log</string>
</dict>
</plist>
